# util/openbsd/endlessh
#!/bin/ksh
#

daemon="/usr/local/bin/endlessh"
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
